# meta/meta-instrumentcluster/recipes-apps/ic-app/files/ic-app.service
[Unit]
Description=Instrument Cluster Application
After=systemd-udev-settle.service weston.service
Requires=systemd-udev-settle.service weston.service
StartLimitInterval=60s
StartLimitBurst=3

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/IC_app
Restart=on-failure
RestartSec=5

# Set proper environment for Wayland client (connects to Weston)
Environment="QT_QPA_PLATFORM=wayland"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="QT_WAYLAND_DISABLE_WINDOWDECORATION=1"
Environment="XDG_RUNTIME_DIR=/run"
Environment="QT_PLUGIN_PATH=/usr/lib/plugins"
Environment="LD_LIBRARY_PATH=/usr/lib:/usr/local/lib"

# Run as root
User=weston

StandardOutput=journal
StandardError=journal
SyslogIdentifier=IC_app

[Install]
WantedBy=multi-user.target
