# meta/meta-headunit/recipes-apps/ambientapp/files/ambientapp.service
[Unit]
Description=Ambient Lighting Application
After=network.target hu-mainapp-compositor.service gearapp.service mediaapp.service
Requires=hu-mainapp-compositor.service
StartLimitInterval=60s
StartLimitBurst=3

[Service]
Type=simple
ExecStartPre=/bin/sleep 7
ExecStart=/usr/bin/AmbientApp
Restart=on-failure
RestartSec=5

# Set proper environment for Qt5 Wayland client (HU_MainApp compositor)
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="QT_QPA_PLATFORM=wayland"
Environment="QT_WAYLAND_DISABLE_WINDOWDECORATION=1"
Environment="WAYLAND_DISPLAY=wayland-1"
Environment="QT_QUICK_BACKEND=software"
Environment="QT_AUTO_SCREEN_SCALE_FACTOR=0"
Environment="QT_SCALE_FACTOR=1"
Environment="QT_FONT_DPI=96"
Environment="QT_PLUGIN_PATH=/usr/lib/plugins:/usr/lib/qt5/plugins"
Environment="QML2_IMPORT_PATH=/usr/lib/qml:/usr/lib/qt5/qml"
Environment="LD_LIBRARY_PATH=/usr/lib:/usr/local/lib"

# Run as root (same as HU_MainApp compositor)
User=weston

StandardOutput=journal
StandardError=journal
SyslogIdentifier=AmbientApp

[Install]
WantedBy=multi-user.target
