# meta/meta-headunit/recipes-core/udev/files/usb-unmount@.service
[Unit]
Description=Unmount USB Drive (%i)

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'umount /media/usb || true'
